# supervisord.conf
[program:redis-server]
command=redis-server

[program:celery]
command={{ PYTHON }} manage.py celery worker --loglevel=info  -Q celery

[program:uwsgi]
command=uwsgi --ini mynewspaper.ini

